# RaspberryPi_OS/pushgateway/pushgateway.service
[Unit]
Description=Push Gateway

[Service]
Restart=on-failure
ExecStart=/usr/local/bin/pushgateway

[Install]
WantedBy=multi-user.target